# crontab
# m h dom mon dow user  command
0 4   * * * /usr/sbin/logrotate /etc/logrotate.conf
